# launchd/com.harness.2026-ai-seamless-gather.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.2026-ai-seamless-gather</string>

  <!-- 파이썬 가상환경 및 수집 스크립트 백그라운드 실행 -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/juntae.park/projects/harness-platform/.venv/bin/python</string>
    <string>/Users/juntae.park/projects/harness-platform/scripts/run_edu_deep_research.py</string>
    <string>--sources</string>
    <string>youtube</string>
    <string>--max-yt-topic-queries</string>
    <string>18</string>
    <string>--max-yt-channel-crawls</string>
    <string>0</string>
    <string>--max-yt-results</string>
    <string>5</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/juntae.park/projects/harness-platform</string>

  <!-- 일 1회 실행으로 축소: YouTube API/수집 비용 억제 우선 -->
  <key>StartInterval</key>
  <integer>86400</integer>

  <!-- 로그 파일 기록 -->
  <key>StandardOutPath</key>
  <string>/Users/juntae.park/projects/harness-platform/logs/2026-ai-seamless-gather.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/juntae.park/projects/harness-platform/logs/2026-ai-seamless-gather.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>ENABLE_YOUTUBE_DATA_API</key>
    <string>false</string>
    <key>YOUTUBE_API_QUERY_LIMIT</key>
    <string>0</string>
    <key>YOUTUBE_API_CHANNEL_LIMIT</key>
    <string>0</string>
    <key>YOUTUBE_TOPIC_QUERY_LIMIT</key>
    <string>18</string>
    <key>YOUTUBE_CHANNEL_CRAWL_LIMIT</key>
    <string>0</string>
  </dict>
</dict>
</plist>
