# launchd/com.harness.daily-ops-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.daily-ops-brief</string>

  <!-- 정밀도 오디트 + 일일 자동화 ops brief 실행 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/juntaepark/projects/harness-platform/scripts/run_daily_ops_brief.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/juntaepark/projects/harness-platform</string>

  <!-- 매일 오전 04:30 Asia/Seoul -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/daily-ops-brief.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/daily-ops-brief.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
