# launchd/com.harness.weekly-business-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.harness.weekly-business-review</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/juntae.park/projects/harness-platform/.venv/bin/python</string>
        <string>/Users/juntae.park/projects/harness-platform/scripts/generate_weekly_business_review.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer> <!-- 5 = Friday -->
        <key>Hour</key>
        <integer>18</integer> <!-- 18:00 KST -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/juntae.park/projects/harness-platform/logs/wbr.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/juntae.park/projects/harness-platform/logs/wbr.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/juntae.park/projects/harness-platform/</string>
</dict>
</plist>
