# harness-os/launchd/com.harness.paper-trading-health-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.paper-trading-health-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/.venv/bin/python</string>
    <string>__ROOT__/scripts/run_paper_trading_health_report.py</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>45</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>__ROOT__/docs/reports/paper_trading_health_report_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/docs/reports/paper_trading_health_report_stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PAPER_ENTRY_LOCK_ALLOWED_SYMBOLS</key>
    <string>ASX,TSM,VRT</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>IBKR_TRADING_MODE</key>
    <string>paper</string>
    <key>PAPER_TRADING_MAX_POSITIONS</key>
    <string>0</string>
    <key>PAPER_PYRAMID_ENABLED</key>
    <string>false</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
