# harness-os/launchd/com.harness.tier3-filter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.tier3-filter</string>

  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/.venv/bin/python3</string>
    <string>__ROOT__/scripts/run_tier3_backlog_worker.py</string>
    <string>--max-batches</string>
    <string>2</string>
    <string>--min-pending</string>
    <string>1</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__ROOT__</string>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__ROOT__/logs/tier3-filter.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/tier3-filter.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>DATABASE_URL</key>
    <string>postgresql://localhost/harness_prod</string>
  </dict>
</dict>
</plist>
