# harness-os/launchd/com.harness.edu-tier3-refine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.edu-tier3-refine</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/juntaepark/projects/harness-platform/.venv/bin/python3</string>
    <string>/Users/juntaepark/projects/harness-platform/scripts/run_edu_tier3_parallel.py</string>
    <string>--workers</string>
    <string>4</string>
    <string>--shard</string>
    <string>0/1</string>
    <string>--min-score</string>
    <string>0.0</string>
    <string>--cost-every</string>
    <string>20</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/juntaepark/projects/harness-platform</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>DATABASE_URL</key>
    <string>postgresql://localhost/harness_prod</string>
  </dict>

  <!-- 매일 10:30 — 10:00 수집/필터 파이프라인 이후 신규 edu를 정제하고 RAG 인덱스까지 증분 갱신 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/edu-tier3-refine.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/edu-tier3-refine.error.log</string>
</dict>
</plist>
