# harness-os/launchd/com.harness.tier2-filter-fast.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.tier2-filter-fast</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/juntaepark/projects/harness-platform/.venv/bin/python3</string>
    <string>/Users/juntaepark/projects/harness-platform/scripts/run_tier2_backlog_worker.py</string>
    <string>--max-batches</string>
    <string>2</string>
    <string>--min-pending</string>
    <string>0</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/juntaepark/projects/harness-platform</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/tier2-filter-fast.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/tier2-filter-fast.error.log</string>
</dict>
</plist>
