# harness-os/launchd/com.harness.pipeline-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harness.pipeline-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/juntaepark/projects/harness-platform/.venv/bin/python3</string>
    <string>/Users/juntaepark/projects/harness-platform/scripts/pipeline_watchdog.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/juntaepark/projects/harness-platform</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>StandardOutPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/pipeline-watchdog.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/juntaepark/projects/harness-platform/logs/pipeline-watchdog.error.log</string>
</dict>
</plist>
